# containerization/files/crontab
PATH=/bin:/usr/bin:/usr/local/bin

*/10 * * * * flock -xn /usr/local/bin/decapod-collect-data timeout --foreground -k 3m 2m /usr/local/bin/decapod-collect-data > /var/log/cron.log 2>&1
