# containerization/files/crontab
0 * * * * /usr/local/bin/cephlcm-cron-clean-expired-tokens > /var/log/cron.log 2>&1
10 3 * * * /usr/local/bin/cephlcm-cron-clean-old-tasks > /var/log/cron.log 2>&1
